# device/rockchip/common/overlays/overlay-yocto/etc/profile.d/weston.sh
#!/bin/sh
export XDG_RUNTIME_DIR=${XDG_RUNTIME_DIR:-/run/user/`id -u weston`}
